# docker/laravel/docker/app/etc/crontab
SHELL=/bin/bash
BASH_ENV=/container.env

* * * * * /usr/local/bin/php /var/www/artisan schedule:run
